# _dev/supervisord.d/celery.ini
[program:celery]
directory=/home/django/xx-scan/
command=/home/django/xs_venv/bin/python xs start celery

process_name=%(program_name)s
numprocs=1
user=root
autostart=true
#autorestart=true
autorestart=unexpected
startsecs=5
startretries=3
priority=2
redirect_stderr=true
stdout_logfile=/var/log/%(program_name)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10

stopasgroup=true
killasgroup=true

